# muse-supervisor/hudson.conf

[program:hudson]
command = java -jar /srv/hudson/hudson.war
process_name = hudson
directory = /srv/hudson
priority = 0
redirect_stderr = true

